# section17/ops/fab_project/roles/web/files/app.conf
# superviserの設定ファイル

# root/workという場所で、guncornというhttpサーバーを用いて、hello.pyのappを実行する
[program:hello]
command = gunicorn --bind 0.0.0.0:5000 hello:app
directory = /root/work/